# systemprompt.service
[Unit]
Description=SystemPrompt Management Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/SystemPrompt-Management/src
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=systemprompt

# Environment
Environment=PYTHONPATH=/home/ec2-user/SystemPrompt-Management/src
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
