# systemd/systemprompt-management.service
[Unit]
Description=SystemPrompt Management
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/SystemPrompt-Management
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
